# deploy/all-seeing-eye.conf
# Upstart script
description     "All Seeing Eye -- A cssauron visualization tool"
author          "Justin Abrahms"

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [06]

# Prepare the environment
#   Create directories for logging and process management
#   Change ownership to the user running the process
pre-start script
    mkdir -p /var/opt/all-seeing-eye
    mkdir -p /var/opt/all-seeing-eye/log
    mkdir -p /var/opt/all-seeing-eye/run
    chown -R all-seeing-eye:all-seeing-eye /var/opt/all-seeing-eye
end script

# If the process quits unexpectadly trigger a respawn
respawn

export PORT=8001

# Start the process
exec start-stop-daemon \
     --start \
     --chuid all-seeing-eye \
     --make-pidfile \
     --pidfile /var/opt/all-seeing-eye/run/all-seeing-eye.pid \
     --exec /usr/bin/npm start /srv/all-seeing-eye/ase/ >> /var/opt/all-seeing-eye/log/all-seeing-eye.log 2>&1
